# systemd/politwoops.target
[Unit]
Wants=politwoops-twitter.service politwoops-worker.service

[Install]
WantedBy=multi-user.target

